# contrib/init/harzcoind.conf
description "Harzcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env HARZCOIND_BIN="/usr/bin/harzcoind"
env HARZCOIND_USER="harzcoin"
env HARZCOIND_GROUP="harzcoin"
env HARZCOIND_PIDDIR="/var/run/harzcoind"
# upstart can't handle variables constructed with other variables
env HARZCOIND_PIDFILE="/var/run/harzcoind/harzcoind.pid"
env HARZCOIND_CONFIGFILE="/etc/harzcoin/harzcoin.conf"
env HARZCOIND_DATADIR="/var/lib/harzcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # harzcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$HARZCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run harzcoind."
        echo "The setting must appear in $HARZCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Harzcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$HARZCOIND_PIDDIR"
    chmod 0755 "$HARZCOIND_PIDDIR"
    chown $HARZCOIND_USER:$HARZCOIND_GROUP "$HARZCOIND_PIDDIR"
    chown $HARZCOIND_USER:$HARZCOIND_GROUP "$HARZCOIND_CONFIGFILE"
    chmod 0660 "$HARZCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$HARZCOIND_PIDFILE" \
    --chuid $HARZCOIND_USER:$HARZCOIND_GROUP \
    --exec "$HARZCOIND_BIN" \
    -- \
    -pid="$HARZCOIND_PIDFILE" \
    -conf="$HARZCOIND_CONFIGFILE" \
    -datadir="$HARZCOIND_DATADIR" \
    -disablewallet \
    -daemon

